# systemd/dmx_audio_react.service

[Unit]
Description=DMX Audio Reactive Light Controller
After=oled_splash.service network-online.target sound.target olad.service
Wants=oled_splash.service network-online.target olad.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pi-dmx-controller
Environment=PYTHONUNBUFFERED=1 ENABLE_TUI=0
ExecStart=/home/pi/pi-dmx-controller/.venv/bin/python /home/pi/pi-dmx-controller/dmx_audio_react.py --no-tui
Restart=on-failure
RestartSec=3
AmbientCapabilities=CAP_SYS_NICE
Nice=-5

[Install]
WantedBy=multi-user.target