# deploy/pi-dmx.service
[Unit]
Description=Pi DMX Controller
After=network-online.target olad.service
Wants=olad.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pi-dmx-controller
Environment=PYTHONUNBUFFERED=1
# Optionally override universe here, e.g.:
# Environment=DMX_UNIVERSE=0
ExecStart=/home/pi/pi-dmx-controller/.venv/bin/python /home/pi/pi-dmx-controller/dmx_audio_react.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
